# install/systemd/pulse-client.service.tmpl
[Unit]
Description=Pulse Client — terminal agent keeping tmux sessions alive
Documentation=https://github.com/kevinzezel/pulse
After=network.target
# Give up after 5 restarts in 60s (prevents tight crash loops from hammering the box).
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=exec
WorkingDirectory=%h/.local/share/pulse/client/src
EnvironmentFile=%h/.config/pulse/client.env
# /snap/bin covers snap-installed VSCode/Cursor (Ubuntu default).
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin:/snap/bin
# tmux attach-session refuses to run without TERM; systemd units inherit no TERM.
Environment=TERM=xterm-256color
# Import the user's graphical session env so the editor-open handler
# (routes/terminal.py:open_editor → subprocess.Popen of VSCode) can reach
# the X server / Wayland compositor / D-Bus user session. Most display
# managers (GDM, SDDM, LightDM with modern setups) run
# `systemctl --user import-environment` on login, populating these in the
# user manager. PassEnvironment= then lets them reach the unit's exec env.
# ExecStartPre is a safety net for distros/setups where the import never
# happened — cheap to re-run, failure is ignored.
PassEnvironment=DISPLAY WAYLAND_DISPLAY XAUTHORITY DBUS_SESSION_BUS_ADDRESS XDG_RUNTIME_DIR XDG_SESSION_TYPE XDG_CURRENT_DESKTOP
ExecStartPre=/bin/sh -c '/usr/bin/systemctl --user import-environment DISPLAY WAYLAND_DISPLAY XAUTHORITY DBUS_SESSION_BUS_ADDRESS XDG_RUNTIME_DIR 2>/dev/null || true'
# We re-source client.env inside the wrapper so the if/else branches see the
# same values that EnvironmentFile= already loaded. systemd's ${VAR} expansion
# can't do conditional arguments — the sh wrapper can.
ExecStart=/bin/sh -c 'set -a; . "$HOME/.config/pulse/client.env"; set +a; if [ "${TLS_ENABLED:-false}" = "true" ]; then exec "$HOME/.local/share/pulse/client/.venv/bin/uvicorn" service:app --host "$API_HOST" --port "$API_PORT" --ssl-keyfile "$TLS_KEY_PATH" --ssl-certfile "$TLS_CERT_PATH"; else exec "$HOME/.local/share/pulse/client/.venv/bin/uvicorn" service:app --host "$API_HOST" --port "$API_PORT"; fi'
Restart=on-failure
RestartSec=5
# Bound the graceful drain of in-flight WebSockets on `systemctl stop` so the
# installer's stop step (used by `pulse upgrade`) cannot block indefinitely if
# uvicorn gets stuck waiting for connections to close. After 15s systemd
# escalates to SIGKILL on the main process and the unit flips to inactive,
# clearing the way for `enable --now` to actually start the new process.
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal

# Kill only the main uvicorn process on stop/restart — not the whole cgroup.
# The client spawns `tmux new-session -d`, which daemonises but stays in this
# unit's cgroup under systemd (fork doesn't escape cgroups v2). With the
# default KillMode=control-group, every `systemctl stop pulse-client` — which
# `pulse upgrade` runs — would take the tmux server down with it, wiping all
# live sessions. `KillMode=process` leaves the tmux daemon running so sessions
# survive upgrades, restarts, and crashes; recover_sessions() reattaches on boot.
KillMode=process

# Security hardening (user unit — most of this is no-op for user units but set for discipline).
# PrivateTmp is intentionally NOT set: the editor-open handler needs access to
# /tmp/.X11-unix/X0 (X server socket) to spawn GUI apps like VSCode, and
# PrivateTmp=true bind-mounts /tmp to a private per-unit dir where that socket
# doesn't exist.
NoNewPrivileges=true

[Install]
WantedBy=default.target
