# install/systemd/pulse-client.service.tmpl
[Unit]
Description=Pulse Client — terminal agent keeping tmux sessions alive
Documentation=https://github.com/kevinzezel/pulse
After=network.target
# Give up after 5 restarts in 60s (prevents tight crash loops from hammering the box).
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=exec
WorkingDirectory=%h/.local/share/pulse/client/src
EnvironmentFile=%h/.config/pulse/client.env
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
# tmux attach-session refuses to run without TERM; systemd units inherit no TERM.
Environment=TERM=xterm-256color
ExecStart=%h/.local/share/pulse/client/.venv/bin/uvicorn service:app --host ${API_HOST} --port ${API_PORT}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Kill only the main uvicorn process on stop/restart — not the whole cgroup.
# The client spawns `tmux new-session -d`, which daemonises but stays in this
# unit's cgroup under systemd (fork doesn't escape cgroups v2). With the
# default KillMode=control-group, every `systemctl stop pulse-client` — which
# `pulse upgrade` runs — would take the tmux server down with it, wiping all
# live sessions. `KillMode=process` leaves the tmux daemon running so sessions
# survive upgrades, restarts, and crashes; recover_sessions() reattaches on boot.
KillMode=process

# Security hardening (user unit — most of this is no-op for user units but set for discipline).
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
