# install/launchd/sh.pulse.dashboard.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>sh.pulse.dashboard</string>
    <!--
      AUTH_PASSWORD / AUTH_JWT_SECRET live in ~/.config/pulse/frontend.env. The
      shell wrapper sources that file so secrets never land in the plist XML —
      avoids XML-escaping issues on weird password characters and lets the user
      rotate AUTH_JWT_SECRET without reinstalling.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <!--
          Next.js (15) refuses to run on Node < 18.18. The launchd PATH is static,
          so if the user has an older /usr/local/bin/node around, npx can pick it
          up even when node@20 is installed via brew. Resolve brew's node@20
          explicitly at start time and prepend it to PATH — works on Apple
          Silicon and Intel without hard-coded paths in the plist.
        -->
        <string>unset NODE_OPTIONS; for brew_bin in /opt/homebrew/bin/brew /usr/local/bin/brew; do [ -x "$brew_bin" ] &amp;&amp; eval "$("$brew_bin" shellenv)" &amp;&amp; break; done; node_prefix="$(brew --prefix node@20 2>/dev/null)"; [ -n "$node_prefix" ] &amp;&amp; PATH="$node_prefix/bin:$PATH"; export PATH; set -a; . "$HOME/.config/pulse/frontend.env"; set +a; exec node server.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>@@HOME@@/.local/share/pulse/frontend</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>@@HOME@@</string>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>NODE_OPTIONS</key>
        <string></string>
        <key>PATH</key>
        <string>@@HOME@@/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <!-- Pin the jsonStore base dir; process.cwd() is unreliable under launchd workers. -->
        <key>PULSE_FRONTEND_ROOT</key>
        <string>@@HOME@@/.local/share/pulse/frontend</string>
        <!-- server.js falls back here for cert/key when frontend.env leaves the paths empty. -->
        <key>PULSE_CONFIG_ROOT</key>
        <string>@@HOME@@/.config/pulse</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <key>StandardOutPath</key>
    <string>@@HOME@@/.local/state/pulse/logs/dashboard.log</string>
    <key>StandardErrorPath</key>
    <string>@@HOME@@/.local/state/pulse/logs/dashboard.err</string>
</dict>
</plist>
